# setting_files/site.conf
[program:site]
command=uwsgi --ini /home/kitakaze/dmoj/site/uwsgi.ini
directory=/home/kitakaze/dmoj/site
stopsignal=QUIT
stdout_logfile=/tmp/site.stdout.log
stderr_logfile=/tmp/site.stderr.log
